# files/etc/systemd/system/debsums.timer
[Unit]
Description=Run debsums daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
