# files/etc/systemd/system/installed-pkgs.timer
[Unit]
Description=Run installed-pkgs daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
