# files/etc/systemd/system/lssys.timer
[Unit]
Description=Run lssys daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
