# files/etc/systemd/system/setugid.timer
[Unit]
Description=Run setugid daily

[Timer]
OnCalendar=daily

[Install]
WantedBy=timers.target
